# openrc/alpine-resize-filesystem
#!/sbin/openrc-run
# Alpine Linux OpenRC init script for alpine-resize-filesystem
# Replaces: debian-resize-filesystem.service / ubuntu-resize-filesystem.service

description="Alpine Linux filesystem resize on first boot"

depend() {
    need localmount
    after alpine-zram-config
    before net
}

start() {
    ebegin "Resizing Alpine root filesystem"
    /usr/lib/alpine/alpine-resize-filesystem start
    eend $?
}
