# openrc/alpine-zram-config
#!/sbin/openrc-run
# Alpine Linux OpenRC init script for alpine-zram-config

description="Alpine Linux ZRAM configuration"

depend() {
    # boot runlevel：在 swap 激活前配置 zram 设备
    need modules
    before swap localmount
}

start() {
    ebegin "Configuring Alpine ZRAM"
    /usr/lib/alpine/alpine-zram-config start
    eend $?
}
