# openrc/alpine-ramlog
#!/sbin/openrc-run
# Alpine Linux OpenRC init script for alpine-ramlog

description="Alpine Linux RAM-based logging"

depend() {
    # boot runlevel：在任何写日志的服务启动前把 /var/log 挂到 tmpfs
    need localmount
    before logger syslog sshd
}

start() {
    ebegin "Starting Alpine RAM logging"
    /usr/lib/alpine/alpine-ramlog start
    eend $?
}

stop() {
    ebegin "Stopping Alpine RAM logging (syncing to disk)"
    /usr/lib/alpine/alpine-ramlog stop
    eend $?
}
